# containers-as-a-service/caas-grafana/scripts/caas-grafana.service
[Unit]
Description=Grafana
After=network-online.target
Wants=network-online.target

[Service]
TimeoutStartSec=15
Type=forking
Restart=on-failure
ExecStop=/opt/service-containers/grafana/scripts/service_stop.sh
ExecStart=/opt/service-containers/grafana/scripts/service_start.sh

[Install]
WantedBy=multi-user.target